# filesystem/etc/rc.local

#!/bin/sh
OPWD=`cd`

# Source any additional environment that was added by the startup script
#. startup.env
mkdir -p /etc/service

# enable services
ln -s /etc/sv/etcd /etc/service/

cd $OPWD